# be-ptconnect/supervisor/laravel-worker.conf
; ============================================
; Supervisor config for PTConnect Queue Worker
; ============================================
; Copy to /etc/supervisor/conf.d/laravel-worker.conf
; Then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start ptconnect-queue:*
; ============================================
; Thay /path/to/your/project bằng đường dẫn thật trên VPS

[program:ptconnect-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/your/project/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/path/to/your/project/storage/logs/queue-worker.log
stopwaitsecs=3600
